# deploy/production/systemd/xblog-api.service
[Unit]
Description=XBlog API
After=network-online.target postgresql.service xblog-minio.service
Wants=network-online.target xblog-minio.service

[Service]
Type=simple
User=xblog
Group=xblog
WorkingDirectory=/srv/xblog/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/xblog/api.env
ExecStart=/usr/bin/env pnpm --filter @xblog/api start
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
