# deploy/production/systemd/xblog-web.service
[Unit]
Description=XBlog public web app
After=network-online.target xblog-api.service
Wants=network-online.target xblog-api.service

[Service]
Type=simple
User=xblog
Group=xblog
WorkingDirectory=/srv/xblog/current
Environment=NODE_ENV=production
EnvironmentFile=/etc/xblog/web.env
ExecStart=/usr/bin/env pnpm --filter @xblog/web start
Restart=always
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
